# infra/launchd/com.budinoff.blissful-cache-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.budinoff.blissful-cache-cleanup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/budinoff/home-lab/Blissful/infra/scripts/blissful-cache-cleanup.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StandardOutPath</key>
  <string>/Users/budinoff/Library/Logs/blissful-cache-cleanup.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/budinoff/Library/Logs/blissful-cache-cleanup.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>BLISSFUL_NAS_ROOT</key>
    <string>/Volumes/2TB/NAS/blissful</string>
  </dict>
</dict>
</plist>
